# services/caduceus-gateway.service
[Unit]
Description=Caduceus Multi-Channel Gateway
After=network.target opencode-serve.service
Wants=opencode-serve.service

[Service]
Type=simple
# CUSTOMIZE: Change WorkingDirectory for your setup
WorkingDirectory=/home/zephyr/projects/galaxy-protocol
ExecStart=/usr/bin/python3 galaxy-protocol/tools/caduceus/gateway.py --config .galaxy/config.json
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
# Optional deterministic OpenCode binary pin for DeepWiki enrichment.
# install-service.sh auto-populates this when it can discover the binary.
Environment=GALAXY_OPENCODE_BIN=

[Install]
WantedBy=multi-user.target
